# debian/dudl-db.cron.monthly
#!/bin/sh


dudl-mus2id
dudl-mus2job >/dev/null

