# debian/dudl-db.cron.daily
#!/bin/sh


dudl-dbavcheck > /dev/null
dudl-dbmaintenance

